# deploy/systemd/supabase-stack.service
[Unit]
Description=Newsticker Supabase Stack
Documentation=https://supabase.com/docs/guides/self-hosting/docker
After=network.target docker.service
Requires=docker.service
PartOf=newsticker.target

[Service]
Type=simple
User=%u
WorkingDirectory=/opt/newsticker/supabase
ExecStart=/usr/bin/docker compose --env-file /opt/newsticker/supabase/.env up
ExecStop=/usr/bin/docker compose --env-file /opt/newsticker/supabase/.env down
Restart=always
RestartSec=10

# Resource limits
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=multi-user.target newsticker.target

