# deploy/host-agent/host-agent.service
[Unit]
Description=Newsticker Host Agent
Documentation=https://github.com/maximiliani/newsticker/tree/main/deploy/host-agent
After=network.target
PartOf=newsticker.target

[Service]
Type=simple
User=%u
Group=%u
WorkingDirectory=/tmp
StandardOutput=journal
StandardError=journal
SyslogIdentifier=host-agent
Environment="PYTHONUNBUFFERED=1"
# Set HOST_AGENT_SECRET if available from environment file, otherwise use a default
EnvironmentFile=-/opt/newsticker/newsticker/.env.local
ExecStart=%python3_path% /opt/newsticker/newsticker/deploy/host-agent/host-agent.py
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=5

# Resource limits
MemoryMax=128M
CPUQuota=50%

[Install]
WantedBy=multi-user.target newsticker.target
